# funambol/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:funambol]
command=/opt/Funambol/bin/funambol start
