# duke/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:tomcat6]
command=/etc/init.d/tomcat6 start
